# VSCodium/User/History/-8ee389c/5RUq.service
[Unit]
Description=Store a log of battery levels