# VSCodium/User/History/-8ee389c/Dkwq.service
[Unit]
Description=Store a log of battery levels

[Service]
Type=oneshot
ExecStart=/bin/bash/var"