# VSCodium/User/History/ef867b4/eHOI.timer
[Unit]
Description=Run battery check every hour

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target